# scripts/kernel-rr.service.tmpl
[Unit]
Description=qbittorrent torrent server

[Service]
User=root
ExecStart=bash {init_path}
Restart=on-abort

[Install]
WantedBy=multi-user.target
